# config/com.aisns.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.aisns.daily</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/butatohitsujitohaiboru/AI-SNS-Automation/scripts/run_daily.sh</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/butatohitsujitohaiboru/AI-SNS-Automation</string>

	<key>RunAtLoad</key>
	<false/>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/butatohitsujitohaiboru/AI-SNS-Automation/logs/launchd.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/butatohitsujitohaiboru/AI-SNS-Automation/logs/launchd.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
</dict>
</plist>
